# provision/hermes-web-ui/systemd/exocortex-hermes-ui.service
[Unit]
Description=Exocortex Hermes Web UI
After=docker.service network-online.target
Requires=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=__REPO_ROOT__/provision/hermes-web-ui
ExecStart=/usr/bin/docker compose --env-file __REPO_ROOT__/provision/hermes-web-ui/env/.env -f __REPO_ROOT__/provision/hermes-web-ui/docker/compose.yml up -d
ExecStop=/usr/bin/docker compose --env-file __REPO_ROOT__/provision/hermes-web-ui/env/.env -f __REPO_ROOT__/provision/hermes-web-ui/docker/compose.yml down
TimeoutStartSec=0

[Install]
WantedBy=multi-user.target
